# package/creality/klipper/files/klipper_service2.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=55
STOP=1
DEPEND=fstab
USE_PROCD=1
PROG=/usr/share/klippy-env/bin/python
USERDATA=/mnt/UDISK/
PRINTER_CONFIG_DIR=$USERDATA/printer_config2
PRINTER_LOGS_DIR=$USERDATA/printer_logs2
PRINTER_CUSTOM_CFG=$PRINTER_CONFIG_DIR/printer.cfg
PRINTER_DEFAULT_CFG=/usr/share/klipper-brain/printer_config2/printer.cfg
VIRTUAL_SD=/root/gcode_files
REAL_SD=$USERDATA/.crealityprint/upload
OOM_ADJ=-17

start_service() {

    [ -d $PRINTER_CONFIG_DIR ] || mkdir -p $PRINTER_CONFIG_DIR
    [ -d $PRINTER_CONFIG_DIR ] && chmod o=rwx $PRINTER_CONFIG_DIR

    [ -s $PRINTER_CUSTOM_CFG ] || return
    [ -s $PRINTER_CUSTOM_CFG ] && chmod o=rwx $PRINTER_CUSTOM_CFG

    [ -d $PRINTER_LOGS_DIR ] || mkdir -p $PRINTER_LOGS_DIR
    [ -d $PRINTER_LOGS_DIR ] && chmod o=rwx $PRINTER_LOGS_DIR

    [ -L $VIRTUAL_SD ] || ln -s $REAL_SD $VIRTUAL_SD
    [ -d $REAL_SD ] && chmod o=rwx $REAL_SD

    procd_open_instance
    procd_set_param env HOME=/root
    procd_set_param oom_adj $OOM_ADJ
    procd_set_param nice -10
    procd_set_param command taskset 7 $PROG /usr/share/klipper/klippy/klippy.py $PRINTER_CUSTOM_CFG -l $PRINTER_LOGS_DIR/klippy.log -a /tmp/klippy_uds2
    procd_close_instance
}
